# backup/rosbackup.cron
0  */6  *  *  *	root  /usr/local/sbin/rosbackup.sh > /dev/null 2>&1
